# docker/crontab
0 0 * * * root /remove_old_data.sh index-1 index-2 >> /var/log/cron.log
